# raspberry_pi/tonex-foot-controller.service
[Unit]
Description=TONEX One Raspberry Pi Web Controller
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi/TonexOneController/raspberry_pi
Environment=PYTHONUNBUFFERED=1
ExecStart=/home/pi/TonexOneController/raspberry_pi/.venv/bin/python /home/pi/TonexOneController/raspberry_pi/web_controller.py --host 0.0.0.0 --port 8080 --tonex-port /dev/ttyACM0
Restart=on-failure
RestartSec=2

[Install]
WantedBy=multi-user.target
